# supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"

[program:nagios3]
command=/bin/bash -c "/usr/sbin/nagios3 -d /etc/nagios3/nagios.cfg"
